# factory/usr/lib/systemd/system/start-measure-boottime.service
[Service]
Type=oneshot
ExecStart=date +%s%N | cut -b1-13 > /tmp/boot-start